# binfrastructure/ros/devel/etc/catkin/profile.d/10.rosconsole.sh
# generated from rosconsole/env-hooks/10.rosconsole.sh.develspace.em

export ROSCONSOLE_CONFIG_FILE=/home/julian/aaad/interactive-manipulation-sandbox/binfrastructure/ros/src/ros_comm/tools/rosconsole/config/rosconsole.config
