# binfrastructure/ros/src/executer_actions/build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.em

export CATKIN_TEST_RESULTS_DIR="/home/julian/aaad/interactive-manipulation-sandbox/binfrastructure/ros/src/executer_actions/build/test_results"
export ROS_TEST_RESULTS_DIR=$CATKIN_TEST_RESULTS_DIR
